# debian/atftpd.socket
[Unit]
Description=Advanced TFTP Server activation socket

[Socket]
ListenDatagram=69

[Install]
WantedBy=sockets.target
