# debian/atftpd.prerm
#!/bin/sh

set -e

# Stops daemon if it is running under our control
if [ -x /etc/init.d/atftpd ]; then 
    if which invoke-rc.d > /dev/null ; then
	invoke-rc.d --quiet atftpd stop
    else
	/etc/init.d/atftpd stop
    fi
fi	

#DEBHELPER#

exit 0
